# .devcontainer/docker_install/vnc/supervisord.conf
[supervisord]
nodaemon=false
user=root
[program:vnc]
command=gosu '$USER' bash '$VNCRUN_PATH'
[program:novnc]
command=gosu '$USER' bash -c "websockify --web=/usr/lib/novnc 80 localhost:5901"